# mmdvm-log-restore.timer
[Unit]
Description=Restore MMDVM log files on startup
 
[Timer]
OnStartupSec=10
 
[Install]
WantedBy=default.target
 
